# com.devbrain.ingest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.devbrain.ingest</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/patrickkelly/devbrain/ingest/.venv/bin/python</string>
        <string>-u</string>
        <string>/Users/patrickkelly/devbrain/ingest/main.py</string>
        <string>watch</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/patrickkelly/devbrain/ingest</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/patrickkelly/devbrain/logs/ingest.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/patrickkelly/devbrain/logs/ingest.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
